# web_service/conf/crontab/web_service-update-crontab
*/1 *   *   *   *   lyc date > /opt/test_web_service/web_service/resources/time.txt
